# misc/etc/rc.local
echo -n 'starting local daemons:'

# Put local daemons here...

echo '.'
